# docker/rootfs/etc/s6/raven/run
#!/bin/sh

cd /usr/local/raven
exec ucode raven.uc
